# Pi System Files/rc.local
#Add this line to the Pi's rc.local file
sudo python3 /home/pi/Scripts/shutdown_pi.py &